# templates/vault-agent/vault-agent.service
[Unit]
Description=Vault Agent for Certificate Renewal
After=network.target

[Service]
ExecStart=/usr/bin/vault agent -config=/etc/vault/vault-agent.hcl
Restart=on-failure
User=vault
Group=vault

[Install]
WantedBy=multi-user.target
